# examples/resticpy.service
[Unit]
Description=Run restic.py
After=plexdrive.service

[Service]
Type=simple
ExecStart=restic.py
