# examples/resticpy.service
[Unit]
Description=Run restic.py

[Service]
Type=simple
User=pi
WorkingDirectory=/home/pi
ExecStart=/home/pi/.local/bin/restic.py